# infra/workspace/codex-app-server.conf
# supervisor program for codex app-server (long-lived JSON-RPC daemon).
# Worker on host connects to ws://<container-ip>:4455 via the per-workspace
# docker bridge — no host-port publish needed.
[program:codex-app-server]
command=codex app-server --listen ws://0.0.0.0:4455
directory=/workspace
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/codex-app-server.stdout.log
stderr_logfile=/tmp/codex-app-server.stderr.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
# Codex reads auth + config from $CODEX_HOME; HOME is already /home/workspace
# via compose env, so ~/.codex resolves to /home/workspace/.codex which is
# backed by a named volume (persistent across container recreates).  The
# auth.json inside is a bind mount from the host user's ~/.codex/auth.json.
environment=HOME="/home/workspace",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
